# mysqld.conf
[program:mysqld_safe]
priority = 100
command = /usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
redirect_stderr = true
stdout_logfile = /var/log/mysqld.supervisor.log
stderr_logfile = /var/log/mysqld.supervisor.error.log
stdout_events_enabled = true
startsecs = 30
autorestart = true
